# crontab
* * * * * python /home/pi/scripts/pillars.py > /dev/null